# roles/offline-backup/templates/offline-backup@.service.j2
[Unit]
Description=Offline Backup for device %i
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/offline-backup run %i
